# traffic.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/traffic

[Install]
WantedBy=multi-user.target 